# rootfs/etc/s6-overlay/s6-rc.d/init-ytdlp2strm/up
/etc/s6-overlay/s6-rc.d/init-ytdlp2strm/run